# deploy/ansible/roles/wpcomposer/templates/wpcomposer.socket.j2
[Unit]
Description=WP Composer Socket

[Socket]
ListenStream={{ go_listen_addr }}
NoDelay=true

[Install]
WantedBy=sockets.target
